# etc/crontab
* * * * * pi /usr/local/rpi-config/bin/publish_temp_metric.sh 2>&1 | logger -t rpitemp
